# cronjobs/crontab
0   */12   *   *   * sudo cp -r /mnt/hd/docker /mnt/usb